# playbooks/roles/data-donation-task/templates/data-donation-task-sync.service.j2
[Unit]
Description=Sync data-donation-task repository
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Remove the existing checkout so we always start from a clean state,
# avoiding stale files, lock conflicts, or partial dependency trees.
ExecStart=/usr/bin/rm -rf {{ data_donation_task_dir }}
ExecStart=/usr/bin/git clone --branch {{ data_donation_task_branch }} {{ data_donation_task_repo }} {{ data_donation_task_dir }}
# --dangerously-allow-all-builds is required because some dependencies (e.g. esbuild,
# better-sqlite3) include native Node add-ons compiled via lifecycle scripts that
# pnpm blocks by default as a security measure.
ExecStart=/usr/bin/pnpm --dir {{ data_donation_task_dir }} install --dangerously-allow-all-builds

[Install]
WantedBy=multi-user.target
